# jasmin/sms_logger.conf
[program:sms_logger]
command=/usr/bin/python /sms_logger.py
autostart=true
autorestart=true
startretries=3
user= jasmin
